# docs/deployment/etc/systemd/system/bemserver-celery.service.d/bemserver-celery.conf
[Service]

# TODO: set actual installation path
Environment="BEMSERVER_CORE_SETTINGS_FILE=/srv/bemserver-core/bemserver-core-settings.py"
Environment="BEMSERVER_CELERY_SETTINGS_FILE=/srv/bemserver-core/bemserver-celery-settings.py"
Environment="CELERY_BIN=/srv/bemserver-core/venv_bemserver_core/bin/celery"
Environment="CELERYD_NODES=w1 w2"
# TODO: ensure queue name matches Celery task_default_queue parameter
Environment="CELERYD_QUEUES='bemserver_core'"
Environment="CELERYD_LOG_LEVEL=WARNING"
